# client.conf
description "vm-announce mDNS client"
author "Jacques Labuschagne <jlabusch@acm.org>"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

respawn

exec /usr/bin/node PLACE/client.js --name HOST
