# server.conf
description "vm-announce mDNS browser"
author "Jacques Labuschagne <jlabusch@acm.org>"

start on runlevel [2345]
stop on runlevel [016]

respawn

exec /usr/bin/node PLACE/server.js
