# upstart.conf
description "vm-announce mDNS client"
author "Jacques Labuschagne <jlabusch@acm.org>"

start on runlevel [345]
stop on runlevel [0126]

respawn

exec /usr/bin/node PLACE/client.js --name HOST
